# systemd/tarsnap-archive-daily.service
[Unit]
Description=Tarsnap Archive Daily

[Service]
Type=simple
ExecStart=/usr/local/bin/tarsnap-archive.sh daily
