# systemd/tarsnap-archive-weekly.service
[Unit]
Description=Tarsnap Archive Weekly

[Service]
Type=simple
ExecStart=/usr/local/bin/tarsnap-archive.sh weekly
